# flm_voice/service/flm-voice.service
[Unit]
Description=flm-voice daemon (hotkey voice-to-text for KDE Plasma Wayland)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Adjust ExecStart to your install location:
#   pipx install -e .   -> %h/.local/bin/flm-voice
#   pip install -e .    inside a venv -> /full/path/to/.venv/bin/flm-voice
ExecStart=%h/.local/bin/flm-voice daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
